# home/.config/autostart/io.elementary.mail.desktop
[Desktop Entry]
Type=Application
Name=io.elementary.mail
X-XDP-Autostart=io.elementary.mail
Exec=io.elementary.mail --background
